# scripts/systemd/oneai-recordings-api.service
[Unit]
Description=OneAI Recordings API
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/oneai-recordings
Environment=NODE_ENV=production
Environment=ONEAI_APP_DIR=/opt/oneai-recordings
ExecStart=/usr/bin/node src/cli.js api
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
