# scripts/systemd/oneai-recordings-scheduler.service
[Unit]
Description=OneAI Recordings Scheduler (legacy Google Meet recorder)
After=network.target oneai-recordings-api.service

[Service]
Type=simple
WorkingDirectory=/opt/oneai-recordings
Environment=NODE_ENV=production
Environment=ONEAI_APP_DIR=/opt/oneai-recordings
ExecStart=/usr/bin/node src/cli.js legacy:scheduler
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
